# .docker/php7/config/crontab
*    *    *    *    *     cd /app/SuiteCRM; php -f cron.php > /dev/null 2>&1
